# systemd/weyriva-waybar.service
[Unit]
Description=Weyriva Waybar panel
PartOf=graphical-session.target
After=graphical-session.target weyriva-ipc.service
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/weyriva component waybar
Restart=on-failure
RestartSec=1
